# scheduling_example/crontab
* * * * * root date >> /usr/src/app/fic_tasks/dates.txt
# New character line required 